# client/etc/crontabs/root
# to test if IPSec connection established every 15 mins
*/15 * * * * /etc/racoon/keep-alive
